# keystone/container-631ac97d8ec2ef5f98e1a354d8209c6b3e004892d68bb086967c0108a9f495e4.service
# container-631ac97d8ec2ef5f98e1a354d8209c6b3e004892d68bb086967c0108a9f495e4.service
# autogenerated by Podman 2.0.5
# Sat Dec 12 11:10:55 CET 2020

[Unit]
Description=Podman container-631ac97d8ec2ef5f98e1a354d8209c6b3e004892d68bb086967c0108a9f495e4.service
Documentation=man:podman-generate-systemd(1)
Wants=network.target
After=network-online.target

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
ExecStart=/usr/bin/podman start 631ac97d8ec2ef5f98e1a354d8209c6b3e004892d68bb086967c0108a9f495e4
ExecStop=/usr/bin/podman stop -t 10 631ac97d8ec2ef5f98e1a354d8209c6b3e004892d68bb086967c0108a9f495e4
ExecStopPost=/usr/bin/podman stop -t 10 631ac97d8ec2ef5f98e1a354d8209c6b3e004892d68bb086967c0108a9f495e4
PIDFile=/run/user/1000/containers/overlay-containers/631ac97d8ec2ef5f98e1a354d8209c6b3e004892d68bb086967c0108a9f495e4/userdata/conmon.pid
KillMode=none
Type=forking

[Install]
WantedBy=multi-user.target default.target
