# keystone/container-859e2af8c6b16b8d72261600198a199e2661b4fe90f687ffd42b1e8ea22a617d.service
# container-859e2af8c6b16b8d72261600198a199e2661b4fe90f687ffd42b1e8ea22a617d.service
# autogenerated by Podman 2.0.5
# Sat Dec 12 11:10:50 CET 2020

[Unit]
Description=Podman container-859e2af8c6b16b8d72261600198a199e2661b4fe90f687ffd42b1e8ea22a617d.service
Documentation=man:podman-generate-systemd(1)
Wants=network.target
After=network-online.target

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
ExecStart=/usr/bin/podman start 859e2af8c6b16b8d72261600198a199e2661b4fe90f687ffd42b1e8ea22a617d
ExecStop=/usr/bin/podman stop -t 10 859e2af8c6b16b8d72261600198a199e2661b4fe90f687ffd42b1e8ea22a617d
ExecStopPost=/usr/bin/podman stop -t 10 859e2af8c6b16b8d72261600198a199e2661b4fe90f687ffd42b1e8ea22a617d
PIDFile=/run/user/1000/containers/overlay-containers/859e2af8c6b16b8d72261600198a199e2661b4fe90f687ffd42b1e8ea22a617d/userdata/conmon.pid
KillMode=none
Type=forking

[Install]
WantedBy=multi-user.target default.target
